# contrib/surface/arm64/usr/lib/systemd/system/pd-mapper.service
# Qualcomm Protection Domain Mapper
# Needed for WiFi (ath12k) and audio DSP firmware bring-up on Snapdragon X Elite.
# Source: https://github.com/andersson/pd-mapper

[Unit]
Description=Qualcomm Protection Domain Mapper
After=qrtr-ns.service
Requires=qrtr-ns.service

[Service]
Type=simple
ExecStart=/usr/libexec/pd-mapper
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
